# certbot.service
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=/bin/certbot renew --post-hook "systemctl restart httpd"
